# init/upstart/kastenwesen_stop.conf
start on stopping docker
task
exec /usr/local/bin/kastenwesen stop