# init/upstart/kastenwesen_start.conf
start on started docker
task
exec /usr/local/bin/kastenwesen start